# pkg/deb/mnt-reform-keyboard-backlight-rainbow/DEBIAN/prerm
#!/bin/sh

# Are we being called to remove the package?
if [ "$1" = "remove" ]; then

  # Stop and disable the service
  systemctl disable --now mnt_reform_keyboard_backlight_rainbow

fi
